# deploy/swarm56-clip.service
[Unit]
Description=Swarm56 클리핑 에이전트 1회 실행 (백오피스 "지금 클리핑" 트리거)
After=network.target

[Service]
Type=oneshot
User=swarm56-web
Group=swarm56-web
WorkingDirectory=/opt/swarm56/agent/current
# 토큰/키는 여기서만 주입(코드/Git 금지). root:swarm56-web 0640.
EnvironmentFile=/etc/swarm56/agent.env
ExecStart=/opt/swarm56/agent/current/.venv/bin/python -m agent.main
# 처리 후 트리거 파일 제거 → 중복 실행 방지(멱등)
ExecStartPost=/bin/rm -f /var/lib/swarm56/triggers/clip.now
TimeoutStartSec=900
NoNewPrivileges=true
PrivateTmp=true
